# deployment/systemd/trading-platform.service
[Unit]
Description=Trading Platform - Professional Trading Application
Documentation=https://github.com/yourusername/trading-platform
After=network.target postgresql.service

# Ensure .env file exists before starting
ConditionPathExists=/var/www/trading-platform/.env
ConditionPathExists=/var/www/trading-platform/dist/index.js

[Service]
Type=simple
User=tradingapp
Group=tradingapp
WorkingDirectory=/var/www/trading-platform

# Environment variables
Environment=NODE_ENV=production
EnvironmentFile=-/var/www/trading-platform/.env

# Preflight validation - check critical environment variables
ExecStartPre=/bin/bash -c 'source /var/www/trading-platform/.env && \
    test -n "$DATABASE_URL" || (echo "ERROR: DATABASE_URL not set" && exit 1) && \
    test -n "$JWT_SECRET" || (echo "ERROR: JWT_SECRET not set" && exit 1) && \
    test -n "$JWT_REFRESH_SECRET" || (echo "ERROR: JWT_REFRESH_SECRET not set" && exit 1)'

# Start command - use env to find node (supports nvm, n, or system node)
ExecStart=/usr/bin/env node /var/www/trading-platform/dist/index.js

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=0

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/trading-platform/logs

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=trading-platform

# Resource limits
LimitNOFILE=65535
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
